# pkill-xwayland.service
[Unit]
Description=Kill Xwayland

[Service]
Type=oneshot
ExecStart=-/usr/bin/pkill -ex -HUP Xwayland
